# templates/systemd.service.j2
# Managed by ansible (role: haproxy)
# DO NOT EDIT!

[Service]
LimitNOFILE={{ haproxy_maxconn + 10 }}
ExecStartPre=
ExecStartPre=/etc/haproxy/gen-conf.sh
ExecStartPre=/usr/sbin/haproxy -f ${CONFIG} -c -q $EXTRAOPTS
ExecReload=
ExecReload=/etc/haproxy/gen-conf.sh
ExecReload=/usr/sbin/haproxy -f ${CONFIG} -c -q $EXTRAOPTS
ExecReload=/bin/kill -USR2 $MAINPID
